# conf.d/chromium.conf
[program:chromium]
user=%(ENV_APP_USERNAME)s
command=sh -c "rm -rf $XDG_CONFIG_HOME/chromium/Singleton* && chromium-browser --start-maximized $CHROMIUM_PARAMS"
autorestart=true
