# conf.d/fluxbox.conf
[program:fluxbox]
user=%(ENV_APP_USERNAME)s
command=fluxbox
autorestart=true
