# conf.d/websockify.conf
[program:websockify]
user=%(ENV_SYS_USERNAME)s
command=websockify --web /usr/share/novnc 5980 localhost:5900 %(ENV_WEBSOCKIFY_PARAMS)s
autorestart=true
