# conf.d/x11vnc.conf
[program:x11vnc]
user=%(ENV_SYS_USERNAME)s
command=x11vnc -forever -shared -nopw
autorestart=true
